# package/org.esxx.ESXX.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
	  	       "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>org.esxx.ESXX</string>

    <key>ProgramArguments</key>
    <array>
      <string>@bindir@/@PACKAGE_NAME@</string>
      <string>--bind</string>
      <string>7654</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>Disabled</key>
    <true/> <!-- Change this to enable the service -->
  </dict>
</plist>
